# debug/.devcontainer/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:rated]
command=/opt/halon/sbin/rated

[program:dlpd]
command=/opt/halon/sbin/dlpd

[program:smtpd]
command=/opt/halon/sbin/smtpd -f